# config/queue-supervisor.conf
; -------------------------------------------------------------------------
; Supervisor config for Laravel queue workers (production)
; Save as: /etc/supervisor/conf.d/cloudbridge-queue.conf
; -------------------------------------------------------------------------

[program:cloudbridge-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/wifi-billing-saas/artisan queue:work database --queue=critical,high,medium,low,default --sleep=1 --tries=3 --timeout=120 --max-time=3600 --max-jobs=1000
directory=/var/www/wifi-billing-saas
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cloudbridge-queue.log
stopwaitsecs=3600

; Optional dedicated critical worker for webhooks and payment callbacks
[program:cloudbridge-queue-critical]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/wifi-billing-saas/artisan queue:work database --queue=critical --sleep=1 --tries=5 --timeout=180 --max-time=3600 --max-jobs=1000
directory=/var/www/wifi-billing-saas
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cloudbridge-queue-critical.log
stopwaitsecs=3600